# modidock/supervisord.conf
# supervisord.conf

[supervisord]
nodaemon=true

[program:mongod]
command=/usr/bin/mongod --bind_ip_all
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/mongod.err.log
stdout_logfile=/var/log/supervisor/mongod.out.log

[program:nodeapp]
command=node index.js
directory=/usr/src/app
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/app.err.log
stdout_logfile=/var/log/supervisor/app.out.log